# roles/ledger-sync/templates/ledger-fetch-all.timer.j2
[Unit]
Description=Run ledger-fetch-all {{ ledger_fetch_schedule }}

[Timer]
OnCalendar={{ ledger_fetch_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
