# roles/ledger-sync/templates/ledger-push.timer.j2
[Unit]
Description=Push ledger data to origin main periodically

[Timer]
OnCalendar={{ ledger_push_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
