# roles/openclaw-scripts/templates/debug-log-forwarder.service.j2
[Unit]
Description=Forward exec tool results to Telegram debug topic

[Service]
Type=oneshot
ExecStart={{ ansible_env.HOME }}/.openclaw/workspace/scripts/debug-log-forwarder.py
