# roles/openclaw-scripts/templates/openclaw-sq-triage.timer.j2
[Unit]
Description=Run sq triage every 2 hours

[Timer]
OnCalendar=0/2:00
Persistent=true

[Install]
WantedBy=timers.target
