# roles/openclaw-scripts/templates/session-sweep.timer.j2
[Unit]
Description=Run session sweep daily at 5am

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true

[Install]
WantedBy=timers.target
