# roles/openclaw-scripts/templates/session-sweep.service.j2
[Unit]
Description=OpenClaw Session Sweep
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 {{ ansible_env.HOME }}/.openclaw/workspace/scripts/session-sweep.py
ExecStart=/bin/bash {{ ansible_env.HOME }}/.openclaw/workspace/scripts/journal-dump-sweep.sh
Environment="PATH={{ ansible_env.HOME }}/.nvm/versions/node/v22.21.1/bin:/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
