# roles/openclaw-scripts/templates/timeline-tick.service.j2
[Unit]
Description=OpenClaw Timeline Tick (per-minute marker reminder)
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 {{ ansible_env.HOME }}/.openclaw/workspace/scripts/timeline-tick.py
Environment="PATH={{ ansible_env.HOME }}/.nvm/versions/node/v22.21.1/bin:/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
